# deploy/supervisor-busao.conf
[program:busao]
command=/home/python/src/busao/bin/start_busao.sh
user=python

